# Paystub.service
[Unit]
Description=Automates Paystub retrieval and Paystub Printing
after=network-online.target

[Service]
Restart=on-failure
WorkingDirectory=/home/vncuser/Documents/Paycheck-Printer
ExecStart=/usr/local/bin/node /home/vncuser/Documents/Paycheck-Printer/Paystub.js

[Install]
WantedBy=multi-user.target
